# packaging/systemd/syswarden.service
# syswarden — system supervision daemon
# Hardened unit per architecture.md §18.
#
# Manual install:
#   sudo install -Dm644 packaging/systemd/syswarden.service \
#                       /etc/systemd/system/syswarden.service
#   sudo systemctl daemon-reload
#   sudo systemctl enable --now syswarden
#
# Manual uninstall:
#   sudo systemctl disable --now syswarden
#   sudo rm /etc/systemd/system/syswarden.service
#   sudo systemctl daemon-reload
#
# Logs:
#   journalctl -u syswarden -f

[Unit]
Description=syswarden system supervision daemon
Documentation=https://github.com/kresley/syswarden
After=local-fs.target
Wants=local-fs.target
# Bound restart storms (architecture.md §18).
StartLimitIntervalSec=300
StartLimitBurst=5

[Service]
Type=notify
ExecStart=/usr/bin/syswarden daemon
# on-failure also covers watchdog timeout, so a hung supervision loop restarts too.
Restart=on-failure
RestartSec=5

# Liveness watchdog: the daemon pings WATCHDOG=1 after each completed tick via
# sd_notify ($NOTIFY_SOCKET, AF_UNIX — architecture.md §17/§18). Keep this well
# above the maximum adaptive poll interval to avoid false positives on long
# healthy sleeps; the daemon caps its own sleep to WatchdogSec/2 as a safety net.
WatchdogSec=180

# State directories — created automatically by systemd (mode 0750, owned by root).
StateDirectory=syswarden
LogsDirectory=syswarden

# Hardening (architecture.md §18)
NoNewPrivileges=yes
ProtectSystem=strict
ProtectHome=yes
PrivateTmp=yes
ProtectKernelTunables=yes
# cgroup access is required to read/govern resource control — must not restrict.
ProtectControlGroups=no
# D-Bus and sd_notify only; no network sockets of any kind.
RestrictAddressFamilies=AF_UNIX
IPAddressDeny=any
SystemCallFilter=@system-service
MemoryDenyWriteExecute=yes
# State files (history, audit, rollback) live here.
ReadWritePaths=/var/lib/syswarden

[Install]
WantedBy=multi-user.target
